# DEBIAN/prerm
#!/bin/sh
set -eu
if [ "$1" = "remove" ] || [ "$1" = "deconfigure" ]; then
  systemctl disable --now barhelper-kiosk.service >/dev/null 2>&1 || true
  systemctl unmask getty@tty1.service >/dev/null 2>&1 || true
  systemctl daemon-reload >/dev/null 2>&1 || true
fi
exit 0
